# 55-c_init_fb1.rules
SUBSYSTEM=="graphics", ACTION=="add", KERNELS=="mesonfb", KERNEL=="fb1", RUN+="c_init_fb1.sh"
