# 55-c_init_fb0.rules
SUBSYSTEM=="graphics", ACTION=="add", KERNELS=="mesonfb", KERNEL=="fb0", RUN+="c_init_fb0.sh"
