# 55-c_init_net.rules
SUBSYSTEM=="net", ACTION=="add", KERNELS=="meson-eth", KERNEL=="eth0", RUN+="c_init_net.sh"
