# 55-c_init_ppmgr.rules
SUBSYSTEM=="ppmgr", ACTION=="add", KERNEL=="ppmgr", RUN+="c_init_ppmgr.sh"
